# assets/base-services/router.slice
[Unit]
Description=Linux Router

[Slice]
