# services/nat-loopback/natloopback.path
[Unit]
Description=watch firewalld policy change by user

[Install]
WantedBy=router.target

[Path]
PathChanged=/etc/firewalld/policies/ipv4-NAT.xml
TriggerLimitBurst=0
