# services/qos/qos.service
[Unit]
Description=QoS setting script for ppp0

[Service]
ExecStart=/usr/bin/bash "${__dirname}/execute.sh"
RemainAfterExit=no
Type=oneshot
