# assets/router.slice
##### ROUTER GENERATED
## source=/data/DevelopmentRoot/github.com/gongt/router.target/assets/router.slice

[Unit]
Description=Linux Router

[Slice]
CPUWeight=100
CPUQuota=60%

