# assets/router.slice
[Unit]
Description=Linux Router

[Install]

[Slice]
CPUWeight=800

