# services/proxy/proxy-reconfigure.service
[Unit]
Description=watch proxy subscription update event

[Service]
Type=oneshot
Restart=no
# todo
# ExecStart=
