# services/miniupnpd/require-nftables-dropin.service
[Unit]
BindsTo=miniupnpd.service
StopWhenUnneeded=yes
