# services/ddns/ddns.timer
[Unit]
Description=run ddns every 1 min

[Install]
WantedBy=router.target

[Timer]
OnUnitInactiveSec=1min
OnBootSec=30s

